# example/example.service
[Unit]
Description=gosdd Example Service

[Service]
Type=simple
ExecStart=%h/example

[Install]
WantedBy=default.target
